# campus-market-arms.service
[Unit]
Description=Campus Market with ARMS Monitoring
After=network.target mysql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/campus_market
Environment=ARMS_APP_NAME=campus_market
Environment=ARMS_REGION_ID=cn-hangzhou
Environment=ARMS_LICENSE_KEY=biimgsqhcm@6ca181feeeac0da
Environment=ARMS_IS_PUBLIC=True
ExecStart=/root/campus_market/start_gunicorn_arms.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
